# systemd/monitor-kuota.service
[Unit]
Description=Monitor Kuota SSH User
After=network.target

[Service]
ExecStart=/usr/local/bin/monitor-kuota
Restart=always

[Install]
WantedBy=multi-user.target
